# PKGBUILDs/msm-firmware-loader/msm-firmware-loader.install
post_install() {
    systemctl enable msm-firmware-loader
    systemctl enable msm-firmware-loader-unpack
}
