# PKGBUILDs/archlinux-msm8916/boot.automount
[Unit]
Description=Automount /boot partition
ConditionPathExists=/boot

[Automount]
Where=/boot
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
