# PKGBUILDs/lk2nd-boot/lk2nd-boot.install
post_install() {
    systemctl enable boot.automount
}
