# PKGBUILDs/lk2nd-boot/userdata-partition.rules
SUBSYSTEM=="block", ACTION=="add|change", ENV{ID_PART_ENTRY_NAME}=="userdata", RUN+="/usr/bin/mount-userdata.sh $env{DEVNAME}"
